# etc/init.d/search-engine
#!/usr/bin/env bash

# chkconfig: 2345 99 01
supervisorctl -s unix:///tmp/supervisord_search_engine.sock ${1} `basename ${0}`
