# etc/analyze/supervisord_report_gen_files/report_gen.ini
[program:report-gen]
command=/usr/bin/report-gen
priority=520
autostart=true
killasgroup=true
stopsignal=KILL
stdout_capture_maxbytes=1MB
redirect_stderr=true
stdout_logfile=/var/log/analyze/report_gen_stdout.log
stderr_logfile=/dev/null
startsecs=5
exitcodes=0                   ; 'expected' exit codes for process (default 0,2)
user=root
